# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   

[supervisord]
logfile=/var/log/supervisord.log 
logfile_maxbytes=50MB        
logfile_backups=10          
loglevel=info              
pidfile=/tmp/supervisord.pid 
nodaemon=true               
childlogdir=/prometheus/log     
umask=022

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:supercronic]
command=/usr/local/bin/supercronic /app/prometheus/cron-prometheus
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=1
redirect_stderr=true


[program:prometheus]
command=/bin/prometheus -config.file=/etc/prometheus/prometheus.yml -storage.local.path=/prometheus -web.console.libraries=/etc/prometheus/console_libraries -web.console.templates=/etc/prometheus/consoles
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=1
redirect_stderr=true
